# services/media-downloader-.service.d/10-hardening.conf
[Service]
MemoryMax=16G
TasksMax=512
LimitNOFILE=8192
UMask=0077
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectKernelLogs=yes
RestrictSUIDSGID=yes
LockPersonality=yes
CapabilityBoundingSet=
AmbientCapabilities=
RuntimeDirectoryMode=0700
Restart=on-failure
RestartSec=10
