# services/chatterbox-tts.service
[Unit]
Description=Chatterbox TTS standard HTTP server
After=network-online.target
Wants=network-online.target

[Service]
User=mdb
Nice=19

WorkingDirectory=/srv/chatterbox

Environment=HOST=127.0.0.1
Environment=PORT=10250
Environment=CUDA_VISIBLE_DEVICES=1
Environment=CHATTERBOX_DEVICE=cuda
Environment=CHATTERBOX_T3_MODEL=v3
Environment=CHATTERBOX_PT_BR_REPO=ResembleAI/Chatterbox-Multilingual-pt-br
Environment=CHATTERBOX_GENERIC_REPO=ResembleAI/chatterbox
Environment=CHATTERBOX_PRELOAD=1

ExecStart=/srv/chatterbox/.venv/bin/python /srv/chatterbox/http_server.py

Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
