# services/omnivoice-tts.service
[Unit]
Description=OmniVoice TTS HTTP server
After=network-online.target
Wants=network-online.target

[Service]
DynamicUser=yes
Nice=19

WorkingDirectory=/srv/OmniVoice

ExecStart=/srv/OmniVoice/venv/bin/python /srv/OmniVoice/http_server.py

Environment=PORT=10440
Environment=TTS_MODEL_ID=/srv/OmniVoice/model
Environment=CUDA_VISIBLE_DEVICES=1
Environment=OMNIVOICE_DEVICE=cuda:0
Environment=TMPDIR=/run/omnivoice-tts
Environment=HF_HOME=/var/cache/omnivoice-tts/huggingface
Environment=TRANSFORMERS_CACHE=/var/cache/omnivoice-tts/huggingface/transformers
Environment=TORCH_HOME=/var/cache/omnivoice-tts/torch

StateDirectory=omnivoice-tts
CacheDirectory=omnivoice-tts
RuntimeDirectory=omnivoice-tts

#Environment=OMNIVOICE_REF_AUDIO=/path/to/reference.wav
#Environment=OMNIVOICE_REF_TEXT=Transcription of the reference audio.

ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/omnivoice-tts /var/cache/omnivoice-tts /run/omnivoice-tts /tmp
PrivateTmp=yes
NoNewPrivileges=yes
RestrictSUIDSGID=yes
PrivateDevices=no

Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
